# conf/goanysync-sync.service.in
[Unit]
Description=Sync @PACKAGE_NAME@ tmpfs volumes to disk

[Service]
Type=oneshot
ExecStart=@bindir@/@PACKAGE_NAME@ sync
