# conf/goanysync.service.in
[Unit]
Description=@PACKAGE_NAME@ directory symlinker and content syncer
Wants=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@bindir@/@PACKAGE_NAME@ start
ExecStop=@bindir@/@PACKAGE_NAME@ stop

[Install]
WantedBy=multi-user.target
